# rc-daemon/com.claude-rc.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.__USER__.claude-rc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-m</string>
        <string>-s</string>
        <string>__CLAUDE_BIN__</string>
        <string>rc</string>
        <string>--remote-control-session-name-prefix</string>
        <string>rc-</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-rc.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-rc.err</string>

    <key>WorkingDirectory</key>
    <string>__WORKING_DIR__</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
